# deploy/groupie.service
[Unit]
Description=Groupie Tracker Go application
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/Groupie-Tracker
# Set environment variables as needed; override in systemd drop-ins if required
Environment=GROUPIE_INSECURE=true
Environment=PORT=8080
ExecStart=/var/www/Groupie-Tracker/groupie
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
